# Documents/GitHub/sigmatic/neutral-strat/setup/sigmatic-trader.service
[Unit]
Description=Sigmatic Neutral Strategy Paper Trader
Documentation=https://github.com/yassineessabar/sigmatic-investment
After=network-online.target
Wants=network-online.target
StartLimitInterval=0

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/sigmatic/neutral-strat
Environment=PYTHONPATH=/home/ubuntu/sigmatic/neutral-strat/src:/home/ubuntu/sigmatic/neutral-strat
Environment=PATH=/home/ubuntu/sigmatic/neutral-strat/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/sigmatic/neutral-strat/venv/bin/python scripts/live_trader.py config/two_week_config.yaml
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/sigmatic/neutral-strat/logs/systemd.log
StandardError=append:/home/ubuntu/sigmatic/neutral-strat/logs/systemd_error.log

# Resource limits
LimitNOFILE=65536
MemoryLimit=1G
CPUQuota=50%

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/sigmatic/neutral-strat/logs
ReadWritePaths=/home/ubuntu/sigmatic/neutral-strat/reports
ReadWritePaths=/home/ubuntu/sigmatic/neutral-strat/exports
PrivateTmp=true

[Install]
WantedBy=multi-user.target